# systemd-user/groundbreaker.service
[Unit]
Description=Background Changer Script

[Service]
ExecStart=/home/sl/.dot/bin/groundbreaker

[Install]
WantedBy=groundbreaker.timer
