# systemd-user/groundbreaker.timer
[Unit]
Description=Background Changer Script

[Timer]
OnCalendar=*:0
Persistent=true

[Install]
WantedBy=timers.target
