# systemd-user/itunes-backup.timer
[Unit]
Description=iTunes Library "hardlink"-backup

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
