# systemd-user/rTiMan-rss.timer
[Unit]
Description=Torrent RSS checker

[Timer]
OnCalendar=*:*:0
Persistent=true

[Install]
WantedBy=timers.target
