# systemd-user/rTiMan-rss.service
[Unit]
Description=Torrent RSS checker

[Service]
Type=oneshot
ExecStart=/home/sl/dev/rTiMan/rss.sh

[Install]
WantedBy=rTiMan-rss.timer
