# systemd-user/statuscenter-daemon.service
[Unit]
Description=StatusCenter block daemon

[Service]
Type=simple
ExecStart=/home/sl/dev/StatusCenter/daemon.sh

[Install]
WantedBy=default.target
